# builder/files/krakend.conf
# upstart script
description "Starts the krankend"

start on runlevel [2345]

stop on runlevel [016]

setuid krakend
setgid nogroup

exec /usr/bin/krakend run -c /etc/krakend/krakend.json
